# armbian-build/userpatches/overlay/etc/systemd/system/hive-firewall.service
[Unit]
Description=Hive Network Node Firewall Setup
Documentation=https://github.com/your-org/hive-network
After=network.target
Before=ssh.service tailscaled.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-firewall.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# 安全设置
User=root
Group=root

# 确保只运行一次
ConditionPathExists=!/var/lib/hive/firewall-configured

# 创建标记文件
ExecStartPost=/bin/bash -c 'mkdir -p /var/lib/hive && touch /var/lib/hive/firewall-configured'

[Install]
WantedBy=multi-user.target